# opsworks_shoryuken/templates/default/shoryuken-workers.conf.erb
# /etc/init/shoryuken.conf - manage a set of Shoryukens

# This example config should work with Ubuntu 12.04+.  It
# allows you to manage multiple Shoryuken instances with
# Upstart, Ubuntu's native service management tool.
#
# See shoryuken.conf for how to manage a single Shoryuken instance.
#
# Use "stop shoryuken-workers-<%= @app_name %>" to stop all Shoryuken instances.
# Use "start shoryuken-workers-<%= @app_name %>" to start all instances.
# Use "restart shoryuken-workers-<%= @app_name %>" to restart all instances.
# Crazy, right?
#

description "manages the set of shoryuken processes"

# This starts upon bootup and stops on shutdown
start on started network-services
stop on stopping network-services

# Set this to the number of Shoryuken processes you want
# to run on this machine, should usually equal the number of cores
env NUM_WORKERS=<%= @num_workers %>

pre-start script
  for i in `seq 1 ${NUM_WORKERS}`
  do
    start shoryuken-<%= @app_name %> index=$i
  done
end script

post-stop script
  for i in `seq 1 ${NUM_WORKERS}`
  do
    stop shoryuken-<%= @app_name %> index=$i
  done
end script